# collector/t1drl-collector.service
[Unit]
Description=T1-DRL submission collector (FastAPI + SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=t1drl
Group=t1drl
WorkingDirectory=/opt/t1drl-collector
EnvironmentFile=/etc/t1drl/env
ExecStart=/opt/t1drl-collector/.venv/bin/python -m uvicorn server:app --host 127.0.0.1 --port 8181
Restart=on-failure
RestartSec=3

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
ReadWritePaths=/var/lib/t1drl

[Install]
WantedBy=multi-user.target
